# KpVotes/systemd/kp-votes.service
[Unit]
Description=KpVotes

[Service]
WorkingDirectory=/srv/kp-votes
ExecStart=/srv/kp-votes/KpVotes

[Install]
WantedBy=default.target
